# hosts/alpha/systemd/transactional-update.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=Mon *-*-* 04:45
RandomizedDelaySec=0
